# roles/voxhelm_deploy/templates/voxhelm.launchd.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>{{ voxhelm_launchd_label | e }}</string>

    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/env</string>
      <string>bash</string>
      <string>-lc</string>
      <string>{{ voxhelm_start_script_path | e }}</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>{{ voxhelm_launchd_path | e }}</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{ voxhelm_app_dir | e }}</string>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
{% if voxhelm_launchd_keepalive | bool %}
    <true/>
{% else %}
    <false/>
{% endif %}

    <key>StandardOutPath</key>
    <string>{{ voxhelm_launchd_stdout_log | e }}</string>
    <key>StandardErrorPath</key>
    <string>{{ voxhelm_launchd_stderr_log | e }}</string>
  </dict>
</plist>

