# roles/zed/templates/zpool-scrub@.service.j2
[Unit]
Description=ZFS scrub for pool %i
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/sbin/zpool scrub %i
