# roles/paperless_deploy/templates/paperless.service.j2
[Unit]
Description=Paperless-ngx Web Server
After=network.target postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User={{ paperless_username }}
Group={{ paperless_group }}
WorkingDirectory={{ paperless_manage_dir }}
EnvironmentFile={{ paperless_env_file }}
Environment="PATH={{ paperless_virtualenv }}/bin:/usr/bin:/bin"
Environment="PYTHONPATH={{ paperless_virtualenv }}/lib/python{{ paperless_python_version }}/site-packages"
Environment="PAPERLESS_NLTK_DIR={{ paperless_nltk_data_path }}"

ExecStart={{ paperless_virtualenv }}/bin/gunicorn -c {{ paperless_site_root }}/gunicorn.conf.py --bind 0.0.0.0:{{ paperless_app_port }} --workers {{ paperless_task_workers }} --worker-class gthread --threads {{ paperless_threads_per_worker }} --worker-tmp-dir /dev/shm --timeout 120 paperless.wsgi:application

Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
ReadWritePaths={{ paperless_external_storage_root }} {{ paperless_site_root }}

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
