# roles/archive_deploy/templates/archive-quote-classifier.service.j2
[Unit]
Description=Archive quote classifier
After=network-online.target archive.service
Wants=network-online.target

[Service]
Type=oneshot
User={{ archive_user }}
Group={{ archive_group }}
WorkingDirectory={{ archive_django_path }}
EnvironmentFile={{ archive_env_path }}
ExecStart={{ archive_venv_bin }}/python manage.py classify_quote_items --limit {{ archive_quote_classifier_limit }} --since-days {{ archive_quote_classifier_since_days }} --timeout {{ archive_quote_classifier_timeout }}
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths={{ archive_home }} {{ archive_site_path }}
