# roles/zfs_syncoid_replication/templates/zfs-syncoid-replication-failure@.service.j2
[Unit]
Description=Alert on ZFS syncoid replication failure (%i)

[Service]
Type=oneshot
ExecStart={{ zfs_syncoid_replication_alert_script_path }} "%i"
