# roles/navidrome_deploy/templates/navidrome-rescan.timer.j2
# Managed by Ansible - Navidrome rescan schedule
[Unit]
Description=Periodic Navidrome library rescan

[Timer]
OnCalendar={{ navidrome_rescan_timer_on_calendar }}
AccuracySec=5m
Persistent=true

[Install]
WantedBy=timers.target
