# roles/dns_deploy/templates/ddns-update.service.j2
[Unit]
Description=Update Gandi DNS for {{ dns_ddns_subdomain }}.{{ dns_ddns_domain }}
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ddns
EnvironmentFile={{ dns_ddns_config_dir }}/gandi.env
ExecStart={{ dns_ddns_script_path }}
StandardOutput=journal
StandardError=journal

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{ dns_ddns_log_dir }}
