# roles/minecraft_java_deploy/templates/minecraft-java.service.j2
[Unit]
Description=Minecraft Java Edition Server
Documentation=https://minecraft.net
After=network.target

[Service]
Type=simple
User={{ minecraft_java_user }}
Group={{ minecraft_java_group }}
WorkingDirectory={{ minecraft_java_server_dir }}
EnvironmentFile={{ minecraft_java_home }}/.env

{% if minecraft_java_loader_mode == 'forge' %}
# Forge runtime wrapper (execs java directly)
ExecStart={{ minecraft_java_forge_wrapper_script_path }}
{% else %}
# Vanilla start command with optimized JVM settings
ExecStart=/usr/bin/java \
    -Xms{{ minecraft_java_memory_min }} \
    -Xmx{{ minecraft_java_memory_max }} \
    {{ minecraft_java_extra_args }} \
    -jar server.jar nogui
{% endif %}

# Graceful shutdown via RCON
ExecStop=/usr/local/bin/mcrcon -H 127.0.0.1 -P {{ minecraft_java_rcon_port }} -p "${RCON_PASSWORD}" "stop"

# Auto-restart on failure
Restart=on-failure
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ minecraft_java_service_name }}

# Resource limits
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
# Note: ProtectHome disabled because server lives in /home/minecraft and needs write access
ProtectHome=false
ReadWritePaths={{ minecraft_java_home }}

[Install]
WantedBy=multi-user.target
