# roles/openclaw_deploy/templates/openclaw.service.j2
# Managed by Ansible (openclaw_deploy). Do not edit on the host.
[Unit]
Description=OpenClaw Gateway (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ openclaw_site_dir }}
ExecStart=/usr/bin/docker compose -f {{ openclaw_site_dir }}/docker-compose.yml up -d --remove-orphans
ExecStop=/usr/bin/docker compose -f {{ openclaw_site_dir }}/docker-compose.yml down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
