# roles/graphyard_deploy/templates/graphyard-web.service.j2
[Unit]
Description=graphyard web application
After=network.target{% if graphyard_runtime_enabled | bool %} {{ graphyard_runtime_service_name }}.service{% endif %}
{% if graphyard_runtime_enabled | bool %}Requires={{ graphyard_runtime_service_name }}.service
{% endif %}

[Service]
Type=simple
User={{ graphyard_user }}
Group={{ graphyard_group }}
WorkingDirectory={{ graphyard_django_path }}
EnvironmentFile={{ graphyard_env_path }}
# SQLite backend: keep write concurrency constrained to avoid lock storms.
ExecStart={{ graphyard_venv_bin }}/granian --workers {{ graphyard_web_workers }} --interface wsgi --host {{ graphyard_web_bind_host }} --port {{ graphyard_web_bind_port }} --access-log config.wsgi_proxy:application
Restart=always
RestartSec={{ graphyard_web_restart_sec }}

[Install]
WantedBy=multi-user.target
