# roles/python_app_systemd/templates/monitor.service.j2
[Unit]
Description={{ service.name }} monitoring agent
Documentation=file://{{ service_home }}/README.md
After=network.target {{ service.name }}.service
Wants=network.target
StartLimitIntervalSec=30
StartLimitBurst=3

[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutStartSec=30
TimeoutStopSec=15
KillSignal=SIGTERM
KillMode=mixed

# Working directory and user/group
WorkingDirectory={{ service_work_dir }}
User={{ service_user }}
Group={{ service_group }}

# Runtime directories
RuntimeDirectory={{ monitoring_service_name }}
RuntimeDirectoryMode=0750
StateDirectory={{ monitoring_service_name }}
StateDirectoryMode=0750
CacheDirectory={{ monitoring_service_name }}
CacheDirectoryMode=0750
LogsDirectory={{ monitoring_service_name }}
LogsDirectoryMode=0750

# Environment variables
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONIOENCODING=utf-8
{% if service.env is defined %}
EnvironmentFile={{ service_home }}/.env
{% endif %}

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
RestrictNamespaces=~CLONE_NEWUSER
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# File system access
ReadWritePaths={{ service_work_dir }}
ReadWritePaths=%S/{{ monitoring_service_name }}
ReadWritePaths=%C/{{ monitoring_service_name }}
ReadWritePaths=%L/{{ monitoring_service_name }}
ReadWritePaths=%t/{{ monitoring_service_name }}
{% if service.monitoring.additional_paths is defined %}
{% for path in service.monitoring.additional_paths %}
ReadWritePaths={{ path }}
{% endfor %}
{% endif %}

# Process limits
LimitNOFILE=1024
LimitNPROC=512
TasksMax=512

# Service execution
ExecStart={{ venv_path }}/bin/start-agent --db %S/{{ monitoring_service_name }}/db.sqlite3{% if monitoring_telegram_enabled | default(false) %} --enable-telegram{% endif %}
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=multi-user.target